# deploy/mws-viewer.service
[Unit]
Description=MWS Viewer Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/apps/mws-viewer
ExecStart=/apps/mws-viewer/venv/bin/python mws_server.py --host 127.0.0.1 --port 8080
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
